# etc/crontab
0,15,30,45 *     * * *     root   chronic /var/www/download.grml.org/grml-mirrors/bin/run_mirmon
17,47 *          * * *     root   chronic /var/www/download.grml.org/grml-mirrors/bin/stamp-master

